# supervisord.conf
[program:feb_erp]
command=python feb_erp.py
autostart=true
autorestart=true

[program:app]
command=python app.py
autostart=true
autorestart=true
